# umonitor.service
[Unit]
Description=Start managing monitor configuration

[Service]
ExecStart=/usr/bin/umonitor --listen

[Install]
WantedBy=graphical.target
